# conf/rc.local
#!/bin/sh -e


sleep 5

cd /home/pi/s63

pm2 start src --name "s63" 
#screen -AdmS s63 -t home
#screen -S s63 -X screen -t init sh -c -c 'pm2 start src; exec bash'


exit 0

